# debian/DEBIAN/postinst
systemctl enable ap2client.service
if [ -d /run/systemd/system ]; then
	systemctl start ap2client.service
fi

setup() {
	if [ -e $2 ]; then
		cp $2 $2.backup
	fi
	cp $1 $2
}

setup /etc/ap2client/dhcpcd.conf /etc/dhcpcd.conf
setup /etc/ap2client/dnsmasq.conf /etc/dnsmasq.conf
setup /etc/ap2client/hostapd.conf /etc/hostapd/hostapd.conf
setup /etc/ap2client/hostapd /etc/default/hostapd