# debian/DEBIAN/postinst
systemctl enable ap2client.service
if [ -d /run/systemd/system ]; then
	systemctl start ap2client.service
fi